# deploy/shepherd.service
[Unit]
Description=Shepherd — mission control for interactive Claude Code
Documentation=https://github.com/erwins-enkel/shepherd
After=network-online.target herdr.service
Wants=network-online.target
# Wants (not Requires): a herdr that fails to start must not cascade-fail Shepherd — the
# server still boots degraded and DIAGNOSE reports `herdr: offline`, which is the whole
# point of the check. After= only orders us behind it when both are enabled.
Wants=herdr.service

[Service]
Type=simple
WorkingDirectory=%h/.shepherd/app
ExecStart=%h/.bun/bin/bun run src/index.ts
Restart=on-failure
RestartSec=2
# 78 (EX_CONFIG) = herdr preflight failed: binary missing on PATH. A human must install
# herdr — restarting won't fix it, so don't crash-loop on this exit code.
RestartPreventExitStatus=78

# herdr lives in ~/.local/bin; bun in ~/.bun/bin — make both resolvable
Environment=PATH=%h/.bun/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=SHEPHERD_PORT=7330
Environment=SHEPHERD_HOST=127.0.0.1
# localhost defaults only. A direct-served HUD (`tailscale serve` on this node) is trusted
# automatically — Shepherd folds the node's own tailnet host into the allowlist at startup.
# A Tailscale Service front (`tailscale serve` registered under `svc:<name>`) proxying to
# this HUD's port is auto-trusted too, discovered the same way. Only a non-Tailscale reverse
# proxy / custom-DNS front (a host that doesn't appear in `tailscale serve status`) needs its
# host added via SHEPHERD_ALLOWED_HOSTS in ~/.shepherd/env.
Environment=SHEPHERD_ALLOWED_HOSTS=localhost,127.0.0.1,::1,[::1]
# optional overrides (token, repo root, etc.) — file may be absent
EnvironmentFile=-%h/.shepherd/env

# app output keeps appending to the existing log file; unit lifecycle (start/stop/restart)
# is still in journald via `journalctl --user -u shepherd`. The file is size-bounded by
# shepherd-logrotate.timer (copytruncate keeps this held fd writing to the truncated inode). #1212
StandardOutput=append:%h/.shepherd/shepherd.log
StandardError=append:%h/.shepherd/shepherd.log

[Install]
WantedBy=default.target
