# deploy/shepherd-logrotate.service
[Unit]
Description=Shepherd — rotate ~/.shepherd/shepherd.log (size-cap via copytruncate)
Documentation=https://github.com/erwins-enkel/shepherd

[Service]
Type=oneshot
# Self-contained rotator — no external `logrotate` binary (it was absent on some hosts → the log
# grew unbounded). provision.ts/update.sh copy deploy/rotate-shepherd-log.sh to ~/.shepherd; systemd
# expands %h here. Run via /bin/sh so the copied script needs no exec bit.
ExecStart=/bin/sh %h/.shepherd/rotate-shepherd-log.sh

# The script uses plain coreutils (cp, gzip, wc, mv, rm) — give the oneshot a normal PATH so it
# resolves them under the minimal PATH a user timer can inherit.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Output goes to journald: `journalctl --user -u shepherd-logrotate`
