# deploy/shepherd-logrotate.service
[Unit]
Description=Shepherd — rotate ~/.shepherd/shepherd.log (size-cap via copytruncate)
Documentation=https://github.com/erwins-enkel/shepherd

[Service]
Type=oneshot
# The config + state live under ~/.shepherd. systemd expands %h here (logrotate does NOT — the
# config's own log path is pre-templated to the real home by provision.ts/update.sh).
ExecStart=logrotate --state %h/.shepherd/logrotate.state %h/.shepherd/shepherd.logrotate

# logrotate commonly lives in /usr/sbin — include sbin (and /usr/local) so the timer-driven run
# resolves it even under a minimal inherited PATH (the backup unit's PATH omits sbin). This list
# MUST stay in lockstep with the detection PATH in install.sh/update.sh and resolveLogrotate's
# candidates in provision.ts — else detection could match a logrotate this run can't exec.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Output goes to journald: `journalctl --user -u shepherd-logrotate`
