# ci/onboarding-harness/shepherd-onboarding.service
[Unit]
Description=Shepherd onboarding regression harness (nightly)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/Work/shepherd
ExecStart=/usr/bin/env bun run onboarding:test
# A oneshot start has no timeout by default — bound the whole nightly run so a
# wedged scenario is killed instead of pinning the host lock indefinitely.
TimeoutStartSec=2h
